# deploy/ai-stock-analyzer-backup.timer
[Unit]
Description=Daily backup for AI Stock Analyzer

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
