# deploy/ai-stock-analyzer-healthcheck.timer
[Unit]
Description=Periodic health check AI Stock Analyzer Bot

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
