# deploy/ai-stock-analyzer-bot.service
[Unit]
Description=AI Stock Analyzer Bot
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=stockbot
Group=stockbot
WorkingDirectory=/var/www/ai-stock-analyzer-bot
EnvironmentFile=/var/www/ai-stock-analyzer-bot/.env
ExecStart=/var/www/ai-stock-analyzer-bot/venv/bin/python run.py
Restart=always
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

