# deploy/ai-stock-analyzer-restart.service
[Unit]
Description=Restart AI Stock Analyzer Bot after failed health check

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart ai-stock-analyzer-bot.service
